# main/goes-platina-mk1/testdata/bird/supervisord.conf
[supervisord]
nodaemon=true

[program:bird]
command=/root/bird -f -c /etc/bird/bird.conf
autorestart=true

