# main/goes-platina-mk1/testdata/frr/ospf/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/lib/frr/zebra -f /etc/frr/frr.conf
autorestart=true

[program:ospf]
command=/usr/lib/frr/ospfd -f /etc/frr/ospfd.conf
autorestart=true

[program:bgp]
command=/usr/lib/frr/bgpd -f /etc/frr/bgpd.conf
autorestart=true

[program:isis]
command=/usr/lib/frr/isisd -f /etc/frr/isisd.conf
autorestart=true
