# docs/examples/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/lib/quagga/zebra -f /etc/quagga/zebra.conf

[program:bgp]
command=/usr/lib/quagga/bgpd -f /etc/quagga/bgpd.conf
