# main/goes-platina-mk1/testdata/gobgp/supervisord.conf
[supervisord]
nodaemon=true

[program:zebra]
command=/usr/lib/quagga/zebra -f /etc/quagga/zebra.conf
autorestart=true

[program:gobgpd]
command=/root/gobgpd -f /etc/gobgp/gobgpd.conf -l info -p 
autorestart=true


