# main/goes-platina-mk1/testdata/net/dhcp/supervisord.R2.conf
[supervisord]
nodaemon=true

[program:dhcpd]
command=/usr/sbin/dhcpd -f -d -cf /root/dhcpd.conf -lf /root/dhcpd.leases 
autorestart=true

